# bitrix/crontab/crontab.cfg
*/20 * * * * /usr/bin/php -f /var/www/podpolkovnyk/data/INOUT/update_sts.php >/dev/null 2>&1
0 */3 * * * /usr/bin/php -f /var/www/podpolkovnyk/data/www/santehsmart.ru/bitrix/php_interface/include/catalog_export/cron_frame.php 21 >/var/www/podpolkovnyk/data/www/santehsmart.ru/bitrix/php_interface/include/catalog_export/logs/21.txt
